# dist/bz.ceh.spacekeeper-save.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Periodic snapshot of the window-to-space layout into history. Saves are
  deduped (no write when nothing changed) and pruned, and the richest
  arrangement is pinned as the high-water snapshot, so a later collapse
  (windows closed, a display disconnected) never clobbers a good layout.

  Install/remove via `make save-agent-install` / `make save-agent-uninstall`.
  Needs Accessibility (and ideally Screen Recording) granted to spacekeeper.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>bz.ceh.spacekeeper-save</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/haynes/bin/spacekeeper</string>
		<string>save</string>
	</array>
	<key>StartInterval</key>
	<integer>180</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/spacekeeper-save.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/spacekeeper-save.log</string>
</dict>
</plist>
